# server/zok-backend.service
[Unit]
Description=ZOK Backend API
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/zok.armyt.biz/server
ExecStart=/home/ubuntu/node-v22.12.0-linux-x64/bin/node index.js
Restart=always
RestartSec=3

# Только порт; остальное (DB_*, JWT_SECRET) — из server/.env
Environment=NODE_ENV=production
Environment=PORT=5002

[Install]
WantedBy=multi-user.target
